# deploy/kua-discord.service
# Bot Discord (fondation) — agent/discord_bot.py. Template — copié dans /etc/systemd/system/
# mais PAS activé par l'agent. Activation (William, sudo) : sudo systemctl enable --now kua-discord
#
# Requiert DISCORD_BOT_TOKEN dans /srv/kua/.env (sinon refuse de démarrer avec un message clair).
# Le bot PROPOSE seulement ; confirmation humaine obligatoire ; allow_auto reste false.
[Unit]
Description=kua-loop-engine Discord bot (propose -> confirme -> thread)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kua-engine
Group=kua-engine
WorkingDirectory=/home/kua-engine/kua-loop-engine
EnvironmentFile=/srv/kua/.env
Environment=PATH=/home/kua-engine/.local/bin:/home/kua-engine/kua-loop-engine/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/kua-engine/kua-loop-engine/.venv/bin/python -m agent.discord_bot
Restart=always
RestartSec=5
# Durcissement (kua-engine, jamais root).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kua-discord

[Install]
WantedBy=multi-user.target
