# deploy/kua-worker.service
# Worker du Runner (boucle claim → claude -p) — doc 05 / doc 06.
# Template — copié dans /etc/systemd/system/ mais PAS activé par l'agent (sudo = bring-live).
# Activation (AVEC William) : sudo systemctl enable --now kua-worker
#
# Le worker est le SEUL invocateur de `claude -p`. Indépendant du gateway (il parle
# directement à la DB) : aucun Requires croisé → un service qui tombe n'entraîne pas les autres.
[Unit]
Description=kua-loop-engine Worker (boucle Runner, claude -p)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kua-engine
Group=kua-engine
WorkingDirectory=/home/kua-engine/kua-loop-engine
EnvironmentFile=/srv/kua/.env
ExecStart=/home/kua-engine/kua-loop-engine/.venv/bin/python -m runner.cli worker
Restart=always
RestartSec=5
# Durcissement (kua-engine, jamais root ; pas d'élévation, FS protégé).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kua-worker

[Install]
WantedBy=multi-user.target
